# image-classification-predict/config/supervisor.conf
[program:app]
command = uwsgi --ini /config/app.ini

[program:nginx]
command = service nginx restart
